# inst/systemd/corteza-matrix.service
[Unit]
Description=corteza Matrix bot (long-poll loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/Rscript -e 'corteza::matrix_run()'
Restart=always
RestartSec=5
StandardOutput=append:%h/.local/state/corteza/matrix.log
StandardError=append:%h/.local/state/corteza/matrix.log

[Install]
WantedBy=default.target
